# linux/supervisor/easyread/201_50_portal-easyread.conf
[program: 201_50_portal-easyread]
command=java -jar -Dapollo.meta=http://172.30.201.60:7080,http://172.30.201.61:7080 portal-easyread.jar
;environment=JAVA_HOME="/usr/java/jdk1.8.0_121/jre",JAVA_BIN="/usr/java/jdk1.8.0_121/jre/bin/"
directory=/app/server/easyread
user=root
priority=10
autostart=true
autorestart=true
startsecs=1
startretries=3
stopasgroup=true
killasgroup=true
stopsignal=KILL
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/app/supervisor_logs/portal-easyread.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
